# platcar.conf
[program:platcar]
command = /opt/platcar/bin/gunicorn -c /opt/platcar-web/gunicorn_config.py platcar_web.wsgi:application
directory = /opt/platcar-web
user = root